# alarm-central-1.4/alarm_authlog.service
[Unit]
Description=Prvide auth.log public keys for alarm-central

[Service]
ExecStart=/usr/bin/alarm_authlog.sh

[Install]
WantedBy=default.target
