# alarm-central-1.4/debian/alarm-central/lib/systemd/system/alarm_cleanup.service
[Unit]
Description=Clean up clutter from client installations

[Service]
User=alarm
ExecStart=/usr/bin/alarm_cleanup.sh

[Install]
WantedBy=default.target
