# alarm-central-1.4/debian/alarm-central/lib/systemd/system/alarmd.service
[Unit]
Description=FastResponseDaemon for CENTRAL

[Service]
User=alarm
ExecStart=/usr/bin/alarmd_central.sh

[Install]
WantedBy=default.target
